# docker/etc/supervisord.d/locationforecastupdate.ini
[program:wdb2ts]
redirect_stderr=true
command=/usr/bin/wget -O /tmp/aaa.xml "http://localhost/metno-wdb2ts/locationforecastupdate"
process_name = locforecastupdate
