# docker/etc/supervisord.d/locationforecastotherupdate.ini
[program:wdb2tsagri]
redirect_stderr=true
command=/usr/bin/wget -O /tmp/aaa.xml "http://localhost/metno-wdb2ts/locationforecastagriupdate"
process_name = locforecastupdateagri
[program:wdb2tsenergy]
redirect_stderr=true
command=/usr/bin/wget -O /tmp/aaa.xml "http://localhost/metno-wdb2ts/locationforecastenergyupdate"
process_name = locforecastupdateenergy
